# deployment/harmonia.service
# Systemd Service File for Harmonia
# Place this in: /etc/systemd/system/harmonia.service

[Unit]
Description=Harmonia Compatibility Analysis Service
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/harmonia

# Use docker-compose to manage the container
ExecStartPre=/usr/bin/docker-compose -f /opt/harmonia/docker-compose.yml down
ExecStart=/usr/bin/docker-compose -f /opt/harmonia/docker-compose.yml up
ExecStop=/usr/bin/docker-compose -f /opt/harmonia/docker-compose.yml down

# Restart policy
Restart=always
RestartSec=10

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=harmonia

[Install]
WantedBy=multi-user.target
